# etc/rpm-find-changes.cron

# Find rpm changes in the /etc tree
0 0 * * * root /usr/bin/rpm-find-changes --exclude-from=/etc/rpm-find-changes/exclude-etc -O /var/cache/rpm-find-changes/etc.txt /etc

